# app-admin/unogenerator_daemon/files/unogenerator.initd
#!/sbin/openrc-run
# Copyright 2021 Gentoo Foundation

description="The python unogenerator module daemon"
pidfile="/var/run/unogenerator.pid"

depend() {
	use net
}

start() {
	ebegin "Starting unogenerator"
	unogenerator_start --instances "${INSTANCES}" --first_port "${FIRST_PORT}"
	eend $?

}
stop() {
	ebegin "Stopping unogenerator"
	unogenerator_stop
	eend $?
}